# scripts/fibreflow-accounting.service
[Unit]
Description=Accounting App (Next.js)
After=network.target

[Service]
Type=simple
User=velo
Group=velo
WorkingDirectory=/home/velo/accounting
ExecStart=/usr/bin/node /home/velo/accounting/node_modules/.bin/next start -p 3101
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
Environment=PORT=3101
EnvironmentFile=/home/velo/accounting/.env.local

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=accounting

# Security
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/velo/accounting

[Install]
WantedBy=multi-user.target
