# scripts/isaflow.service
[Unit]
Description=IsaFlow App (Next.js)
After=network.target

[Service]
Type=simple
User=velo
Group=velo
WorkingDirectory=/home/velo/isaflow
ExecStart=/usr/bin/node /home/velo/isaflow/server.js
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
Environment=PORT=3101
EnvironmentFile=/home/velo/isaflow/.env.local

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=isaflow

# Security
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/velo/isaflow /tmp

[Install]
WantedBy=multi-user.target
